# files/etc/init.d/myproxy
#!/bin/sh /etc/rc.common

START=99

start() {
    /usr/bin/myproxy_apply.sh
}

restart() {
    start
}